# dnd-display.service
[Unit]
Description=DnD Table Display Server
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=dnd
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/dnd/.Xauthority
WorkingDirectory=/home/dnd/dnd-display
ExecStartPre=/usr/bin/xrandr --display :0 --output HDMI-1 --mode 1920x1080 --rate 60
ExecStartPre=/usr/bin/xset -display :0 -dpms
ExecStartPre=/usr/bin/xset -display :0 s off
ExecStartPre=/usr/bin/xset -display :0 s noblank
ExecStart=/usr/bin/python3 /home/dnd/dnd-display/main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
